# ToUpload/app-supervisor.conf
[program: cp3APIapp]
command=/usr/local/bin/app